# etc/rc.local
#!/bin/sh -e 

# Wireless power saving mode disable

# Legacy method
#/sbin/iwconfig wlan0 power off

# new method 
/usr/sbin/iw wlan0 set power_save off


exit $?
